# files/files/system/shared/etc/profile.d/wine.sh
#!/usr/bin/env bash
export WINEESYNC=1
export WINEFSYNC=1
