# config/files/shared/lib/systemd/system/rtkit-daemon.service.d/limits.conf
[Service]
EnvironmentFile=/etc/sysconfig/rtkit
ExecStart=
ExecStart=/usr/libexec/rtkit-daemon $RTKIT_ARGS
