# config/files/plasma/etc/profile.d/kde.sh
#!/bin/sh
export MOZ_ENABLE_WAYLAND=1
export MOZ_USE_XINPUT2=1
